# saltstack/salt/webapp/gowebapp.conf
# Upstart conf of gowebapp

description "Gowebapp service"
start on (started networking)

respawn
respawn limit 10 5

chdir /opt/gowebapp

exec /opt/gowebapp/gowebapp > /var/log/gowebapp.log 2>&1